# rpm/root/etc/systemd/system/MyWeb.service
[Unit]
Description=My Web

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/usr/bin/bin/env dotnet /usr/shared/my-web/MyWeb.dll